# app.service
[Unit]
Description=<nombre de la aplicaci贸n>

[Service]
ExecStart=/usr/bin/nodemon /var/www/<ubicaci贸n de app.js | index.js>
Restart=always
RestartSec=10
User=nobody
Group=nogroup
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=<development | production> PORT=<puerto de la aplicaci贸n>
WorkingDirectory=/var/www/<ubicaci贸n del proyecto>

[Install]
WantedBy=multi-user.target
